# NetworkManager/dispatcher.d/20-dropboxd
#!/bin/bash

case "$2" in
  up) # $INTERFACE is up
    
    # if pidof dropbox; then
    #   killall dropbox
    # fi
    
    # USERS=$(w | awk '/\/X/ {print $1}')

    # for user in $USERS; do
    #   su -c "DISPLAY=$(w | awk '/\/X/ {print $11}') dropboxd &" $user
    # done
    su -c "DISPLAY=:0 dropboxd &" orschiro &
  ;;
  down) # $INTERFACE is down
    killall dropbox
  ;;
esac