# rc.local
tee <<< disable  /proc/acpi/ibm/bluetooth
