# systemd/system/insync@.service
.include /usr/lib/systemd/system/insync@.service
[Service]
KillMode=control-group
