# systemd/system/lirc-logitech-r400.service
[Service]
ExecStart=/usr/bin/lircd --nodaemon --driver=devinput --device=/dev/input/event17

