# systemd/system/lirc.service.d/logitech-r400.conf
[Service]
ExecStart=
ExecStart=/usr/bin/lircd --driver=devinput --device=/dev/input/ir

